# examples/systemd-user/etesync-dav.service
[Unit]
Description=Cal/CardDAV frontend for Etesync

[Service]
Type=simple
ExecStart=/usr/bin/etesync-dav

[Install]
WantedBy=default.target
